# _superninja_startup.conf

[program:18080_psiphon-tunnel-]
command=/bin/bash -c "source ~/.bashrc 2>/dev/null; exec /opt/ninja-browser-use/psiphon-tunnel-core-x86_64 -config /opt/ninja-browser-use/psiphon.config.json"
directory=/opt/ninja-browser-use
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/18080_psiphon-tunnel-.err.log
stdout_logfile=/var/log/supervisor/18080_psiphon-tunnel-.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2

